# systemd/log_eraser.service
[Unit]
   Description=Periodic Log Eraser

[Service]
   Type=oneshot
   ExecStart=/usr/bin/erase_log.sh
