# systemd/log_eraser.timer
[Unit]
   Description=Periodic Log Eraser

[Timer]
   OnCalendar=daily

[Install]
   WantedBy=timers.target
